# contrib/init/peracoind.conf
description "Peracoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env PERACOIND_BIN="/usr/bin/peracoind"
env PERACOIND_USER="peracoin"
env PERACOIND_GROUP="peracoin"
env PERACOIND_PIDDIR="/var/run/peracoind"
# upstart can't handle variables constructed with other variables
env PERACOIND_PIDFILE="/var/run/peracoind/peracoind.pid"
env PERACOIND_CONFIGFILE="/etc/peracoin/peracoin.conf"
env PERACOIND_DATADIR="/var/lib/peracoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # peracoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$PERACOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run peracoind."
        echo "The setting must appear in $PERACOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Peracoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$PERACOIND_PIDDIR"
    chmod 0755 "$PERACOIND_PIDDIR"
    chown $PERACOIND_USER:$PERACOIND_GROUP "$PERACOIND_PIDDIR"
    chown $PERACOIND_USER:$PERACOIND_GROUP "$PERACOIND_CONFIGFILE"
    chmod 0660 "$PERACOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$PERACOIND_PIDFILE" \
    --chuid $PERACOIND_USER:$PERACOIND_GROUP \
    --exec "$PERACOIND_BIN" \
    -- \
    -pid="$PERACOIND_PIDFILE" \
    -conf="$PERACOIND_CONFIGFILE" \
    -datadir="$PERACOIND_DATADIR" \
    -disablewallet \
    -daemon

